# conf/sinkhole/sinkhole.cron.d
#*/1 * * * * root /home/bind_sinkhole/export_host.sh > /etc/bind/named.conf.sinkhole.swp
*/1 * * * * root /home/bind_sinkhole/export_host.sh 
*/1 * * * * root /home/bind_sinkhole/reload_bind.sh 
